# your-service-name.service
[Unit]
Description=My Awesome Application
After=network.target

[Service]
Type=simple
User=yourusername
WorkingDirectory=/path/to/your/application
ExecStart=/usr/bin/java -jar your-application.jar
Restart=on-failure

[Install]
WantedBy=multi-user.target
